# packaging/arch/wepapered.install
post_install() {
	cat <<'EOF'
:: wepapered installed.
   Start the daemon (per user):  systemctl --user enable --now wepapered
   or run it directly:           wepapered-daemon
   Browse wallpapers:            wepaperedctl gui   (or the "Wepapered" launcher entry)
   Settings:                     wepaperedctl settings
   Needs Wallpaper Engine (Steam under Proton) and Hyprland running.
EOF
}

post_upgrade() {
	post_install
}
